# roles/htpc/files/flatpak-update.service
[Unit]
Description=Flatpak update

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak update -y